# packaging/bluetooth-map-agent.service
[Unit]
Description=Bluetooth MAP Agent

[Service]
Type=dbus
BusName=org.bluez.map_agent
ExecStart=/usr/bin/bluetooth-map-agent
